# crontab
MAILTO=""
0 12 * * * push-dotfiles >> /var/log/push-dotfiles.log 2>&1
